# http/template.socket
[Unit]
Description=Socket for example

[Socket]
ListenStream=8080
ReusePort=true
SocketMode=0660

[Install]
WantedBy=sockets.target
